# deploy/blame-the-guilty-tunnel.service
[Unit]
Description=Blame the Guilty ngrok tunnel
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ngrok http --domain=moonlike-silenced-sprung.ngrok-free.dev 5000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
